# lottery.service
[Unit]
Description=Lottery Web Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/nginx/html
ExecStart=/usr/bin/python3 /usr/local/nginx/html/app.py
Restart=always
RestartSec=10
StandardOutput=append:/usr/local/nginx/html/lottery.log
StandardError=append:/usr/local/nginx/html/lottery_error.log

# 环境变量
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
